# systemd/openclaw-reload.service
[Unit]
Description=Auto daemon-reload on openclaw unit change

[Service]
Type=oneshot
ExecStart=/bin/systemctl daemon-reload
ExecStartPost=/bin/systemctl restart openclaw.service
