# systemd/openclaw-reload.path
[Unit]
Description=Watch for openclaw.service unit file changes

[Path]
PathChanged=/etc/systemd/system/openclaw.service
PathChanged=/etc/systemd/system/openclaw.service.d/override.conf

[Install]
WantedBy=multi-user.target
